# launchd/com.dragon.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dragon.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-B</string>
        <string>/Users/ashish/Documents/beast-trader/scripts/watchdog.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ashish/Documents/beast-trader</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/ashish/Documents/beast-trader/logs/watchdog_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ashish/Documents/beast-trader/logs/watchdog_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONDONTWRITEBYTECODE</key>
        <string>1</string>
        <key>HOME</key>
        <string>/Users/ashish</string>
    </dict>
</dict>
</plist>
